# roles/with/home/mounts/files/media-nfs-home-homes.automount
[Unit]
Description=Automount NAS homes share using NFS

[Automount]
Where=/media/nfs/home/homes
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
