# roles/with/home/mounts/files/media-nfs-home-music.automount
[Unit]
Description=Automount NAS music share using NFS

[Automount]
Where=/media/nfs/home/music
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
