# roles/home/mounts/files/media-nfs-home-doc.automount
[Unit]
Description=Automount NAS doc share using NFS

[Automount]
Where=/media/nfs/home/doc
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
