# roles/system/install/files/efi.automount
[Unit]
Description=Automount EFI

[Automount]
Where=/efi
TimeoutIdleSec=10min

[Install]
WantedBy=multi-user.target
